# internal/systemd/templates/llmem-dream.service
[Unit]
Description=LLMem Dream Consolidation
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/llmem dream --apply
WorkingDirectory={{.HomeDir}}
Environment=LMEM_HOME={{.HomeDir}}

[Install]
WantedBy=multi-user.target