# internal/systemd/templates/llmem-dream.timer
[Unit]
Description=LLMem Dream Consolidation Timer

[Timer]
OnCalendar={{.DreamSchedule}}
Persistent=true

[Install]
WantedBy=timers.target